# etc/profile.d/00_bash.sh
alias ll="ls -lh"
alias la="ls -lah"
alias ..="cd .."

