# etc/profile.d/10_base.sh

export EPICS_BASE=`echo /ics/tools/base-*`

export PATH="$EPICS_BASE/bin/linux-x86_64:$PATH"

