# etc/profile.d/10_base.sh

export EPICS_BASE=`echo /ics/tools/base-*`
export EPICS_HOST_ARCH=linux-x86_64

export PATH="$EPICS_BASE/bin/$EPICS_HOST_ARCH:$PATH"

