# etc/profile.d/20_java.sh

export JAVA_HOME=/ics/tools/jdk
export PATH="$JAVA_HOME/bin:$PATH"
