# etc/profile.d/30_p4p.sh
export PYTHONPATH=/ics/tools/p4p/python3.8/linux-x86_64
